# scripts/systemd/intellitask-updater.timer
[Unit]
Description=Run IntelliTask auto-updater every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
