# scripts/systemd/intellitask-updater.service
[Unit]
Description=IntelliTask Auto Updater Service
Wants=intellitask-updater.timer

[Service]
Type=oneshot
EnvironmentFile=/home/ubuntu/app/.env
ExecStart=/usr/local/bin/intellitask-update.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
